# deploy/review-bot.service
[Unit]
Description=Review Bot – Telegram notifications for restaurant reviews
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=reviewbot
WorkingDirectory=/opt/review-bot
EnvironmentFile=/opt/review-bot/.env
ExecStart=python -m bot.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
